# cron/debian/wwmap-cron.cron.d
0 21 * * * wwmap wwmap-notifier
0 0 * * 7 wwmap wwmap-backup
12 * * * * wwmap wwmap-catalog-sync
