# cron/debian/wwmap-cron.cron.d
0 21 * * * wwmap wwmap-notifier
0 0 * * 7 wwmap wwmap-backup
12 * * * * wwmap wwmap-catalog-sync
0 0 1 * * wwmap wwmap-db-clean
0 23 * * * wwmap wwmap-spot-sort
